# config/service/memcached/run
#!/bin/bash

exec /usr/bin/memcached start -u memcache > /dev/null 2>&1 & wait